# php/supervisor/conf.d/php.conf
[supervisord]
nodaemon=true

[program:php]
command = docker-php-entrypoint php-fpm
